# backend/.deploy/supervisor/supervisor.ini
[supervisord]
nodaemon=true
logfile=/app/storage/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777
chown=www-data:www-data

[program:php-fpm]
directory=/app
command=php-fpm
priority=4
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --queue=queue --sleep=3 --tries=5 --max-time=3600
directory=/app
priority=5
numprocs=2
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/queue.log

[program:laravel-shedule]
process_name=%(program_name)s_%(process_num)02d
command=php artisan schedule:work
directory=/app
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/app/storage/logs/worker.log
stderr_logfile=/app/storage/logs/worker.log.error
stopwaitsecs=3600
